# linux/blackout.cron
# Update the blocklist every hour.
0 * * * * root /usr/local/sbin/blackout.sh
